# ndjbdns/etc/supervisord.d/dnscache-int.ini
[program:dnscache-int]
command = /usr/sbin/tinydns
directory = /var/service/dnscache-int
stdout_logfile = /var/log/supervisor/dnscache-int.log
stderr_logfile = /var/log/supervisor/dnscache-int.err
