# ndjbdns/etc/supervisord.d/tinydns-ext.ini
[program:tinydns-ext]
command = /usr/sbin/tinydns
directory = /var/service/tinydns-ext
stdout_logfile = /var/log/supervisor/tinydns-ext.log
stderr_logfile = /var/log/supervisor/tinydns-ext.err
