# djbdns/etc/supervisord.d/dnscache-ext.ini
[program:dnscache-ext]
command = /opt/djbtools/bin/tinydns
directory = /var/service/dnscache-ext
stdout_logfile = /var/log/supervisor/dnscache-ext.log
stderr_logfile = /var/log/supervisor/dnscache-ext.err
