# djbdns/etc/supervisord.d/dnscache-int.ini
[program:dnscache-int]
command = /opt/djbtools/bin/tinydns
directory = /var/service/dnscache-int
stdout_logfile = /var/log/supervisor/dnscache-int.log
stderr_logfile = /var/log/supervisor/dnscache-int.err
