# djbdns/etc/supervisord.d/tinydns-int.ini
[program:tinydns-int]
command = /opt/djbtools/bin/tinydns
directory = /var/service/tinydns-int
stdout_logfile = /var/log/supervisor/tinydns-int.log
stderr_logfile = /var/log/supervisor/tinydns-int.err
